# supervisor/conf.d/mongod.conf
[program:mongod]
command=/usr/local/bin/mongod
autostart=false
autorestart=true
