# supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/local/opt/nginx/sbin/nginx -g "daemon off;"
autostart=false
autorestart=true
